# scripts/supervisor.conf
; ============================================================
; Aviva HealthCare — Supervisor Queue Worker Config
; ============================================================
; Installation:
;   sudo apt install supervisor
;   sudo cp scripts/supervisor.conf /etc/supervisor/conf.d/clinic.conf
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start clinic-worker:*
;
; Check status:
;   sudo supervisorctl status
; ============================================================

[program:clinic-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/clinic/artisan queue:work database --sleep=3 --tries=3 --timeout=120 --max-jobs=500 --max-time=3600
directory=/var/www/clinic
autostart=true
autorestart=true
stopwaitsecs=120
numprocs=2
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/clinic/storage/logs/worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5

[program:clinic-scheduler]
; Runs Laravel scheduler every minute (alternative to cron)
command=php /var/www/clinic/artisan schedule:work
directory=/var/www/clinic
autostart=true
autorestart=true
stopwaitsecs=30
numprocs=1
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/clinic/storage/logs/scheduler.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
